# household-app.service
[Unit]
Description=Household Management Application
Requires=docker.service
After=docker.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/household-app
ExecStart=/usr/local/bin/docker-compose --env-file .env.production up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=300
User=root
Group=docker

# Restart policy
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
